# etherpad-lite/etherpad.conf
[program:etherpad]
directory=/opt/etherpad
command=/bin/bash ./start.sh
user=root
autostart=true
autorestart=true
